# res/DEBIAN/postinst
#!/bin/bash

set -e

if [ "$1" = configure ]; then

	INITSYS=$(ls -al /proc/1/exe | awk -F' ' '{print $NF}' | awk -F'/' '{print $NF}')
    # Create symlink for Xfinity binary
    if [ ! -e /usr/bin/xfinity ]; then
        ln -s /usr/lib/xfinity/xfinity /usr/bin/xfinity
    fi

	if [ "systemd" == "$INITSYS" ]; then
 		# Remove any conflicting service files
        /bin/rm -f /etc/systemd/system/xfinity.service /usr/lib/systemd/system/xfinity.service /usr/lib/systemd/user/xfinity.service >/dev/null 2>&1
		version=$(python3 -V 2>&1 | grep -Po '(?<=Python )(.+)')
		parsedVersion=$(echo "${version//./}")
        mkdir -p /usr/lib/systemd/system/
        cp /usr/share/xfinity/files/systemd/xfinity.service /usr/lib/systemd/system/xfinity.service

		# Fix potential path issue in service file
        if [ -e /usr/bin/pkill ]; then
            /bin/sed -i "s|pkill|/usr/bin/pkill|g" /usr/lib/systemd/system/xfinity.service
        fi

        # Reload and enable the service
        systemctl daemon-reload
        systemctl enable xfinity
        systemctl start xfinity || {
            echo "Failed to start xfinity service"
            exit 1
        }
	fi
fi
