# res/DEBIAN/prerm
#!/bin/bash

set -e

case $1 in
    remove|upgrade)
		INITSYS=$(ls -al /proc/1/exe | awk -F' ' '{print $NF}' | awk -F'/' '{print $NF}')
		rm -f /usr/bin/xfinity
		
		if [ "systemd" == "${INITSYS}" ]; then

			systemctl stop xfinity || true
			systemctl disable xfinity || true
			rm /etc/systemd/system/xfinity.service /usr/lib/systemd/system/xfinity.service || true

			# workaround temp dev build between 1.1.9 and 1.2.0
			serverUser=$(ps -ef | grep -E 'xfinity +--server' | grep -v 'sudo ' | awk '{print $1}' | head -1)
			if [ "$serverUser" != "" ] && [ "$serverUser" != "root" ]
			then
				systemctl --machine=${serverUser}@.host --user stop xfinity >/dev/null 2>&1 || true
			fi
			rm /usr/lib/systemd/user/xfinity.service >/dev/null 2>&1 || true
		fi
        ;;
esac

exit 0
